# oscserial.service
[Unit]
Description=oscserial
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/sousaFX-rpi-scripts/oscserial.py
WorkingDirectory=/home/pi/sousaFX-rpi-scripts/
Restart=always
User=root

[Install]
WantedBy=multi-user.target
